# kafka/supervisor/kafka.conf
[program:kafka]
command=/kafka/bin/kafka-server-start.sh /kafka/config/server.properties
autostart=true
autorestart=true