# kafka/supervisor/zookeeper.conf
[program:zookeeper]
command=/kafka/bin/zookeeper-server-start.sh /kafka/config/zookeeper.properties
autostart=true
autorestart=true