# kafka/supervisor/zookeeper.conf
[program:zookeeper]
command=/kafka/bin/zookeeper-server-start.sh /kafka/config/zookeeper.properties 2>&1 > /dev/null
autostart=true
autorestart=true